# CrewClaw/daemons/config/launchd/lucas-daemon.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>ai.homeai.daemon</string>
    <key>ProgramArguments</key>
    <array>
        <string>/opt/homebrew/bin/node</string>
        <string>/Users/xinbinanshan/HomeAI/homeai/homeai-daemon/index.js</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/xinbinanshan/HomeAI</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>NODE_ENV</key>
        <string>production</string>
        <key>HOMEAI_PORT</key>
        <string>3000</string>
    </dict>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
        <key>CrashCountThreshold</key>
        <string>5</string>
    </dict>
    <key>StandardOutPath</key>
    <string>/Users/xinbinanshan/HomeAI/Logs/homeai-daemon.stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/xinbinanshan/HomeAI/Logs/homeai-daemon.stderr.log</string>
    <key>ProcessType</key>
    <string>Interactive</string>
    <key>ThrottleInterval</key>
    <integer>10</integer>
</dict>
</plist>
