# daemons/config/launchd/lisa-daemon.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>ai.homeai.lisa</string>
    <key>ProgramArguments</key>
    <array>
        <string>/opt/homebrew/bin/node</string>
        <string>/Users/xinbinanshan/HomeAI/homeai/lisa-daemon/index.js</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/xinbinanshan/HomeAI</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>NODE_ENV</key>
        <string>production</string>
        <key>LISA_PORT</key>
        <string>3002</string>
    </dict>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
        <key>CrashCountThreshold</key>
        <string>5</string>
    </dict>
    <key>StandardOutPath</key>
    <string>/Users/xinbinanshan/HomeAI/logs/lisa-daemon.stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/xinbinanshan/HomeAI/logs/lisa-daemon.stderr.log</string>
    <key>ProcessType</key>
    <string>Interactive</string>
    <key>ThrottleInterval</key>
    <integer>10</integer>
    <key>DependentParent</key>
    <string>ai.homeai.andy</string>
</dict>
</plist>
